# roles/bootstrap/templates/etc/local.d/00-user-pi.start
/usr/bin/id pi 2>&1 >/dev/null
if [ $? -ne 0 ]; then
   addgroup pi
   adduser -D -h /home/pi -G pi -s /bin/ash pi
   addgroup pi wheel
   passwd -d pi
fi
install -d -m 700 -o pi -g pi /home/pi
