# ansible/roles/radiogaga/templates/etc/init.d/radiogagad
#!/sbin/openrc-run

command=/usr/local/bin/radiogagad
command_background=true
command_args_backgrounds="--daemon"
command_user="radiogagad:gpio"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/radiogaga.log"

start_pre() {
   touch ${output_log} && chown ${command_user} ${output_log}
}
