# roles/radiogaga/templates/etc/local.d/90-download-playlist.start
python /etc/local.d/download-playlists.py
