# roles/radiogaga/templates/etc/local.d/91-mpc-start-playlist.start
if [ -f /tmp/playlist.start ]; then
    pl=$(head -1 /tmp/playlist.start)
    while [ -z "$(mpc current)" ]; do
        mpc load "$pl"
        mpc play
        sleep 1
    done
fi
